# deploy/systemd/work-time-cloudflared.service
[Unit]
Description=Work Time Cloudflare Quick Tunnel + KV updater
After=network-online.target work-time-api.service
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/srv/app
EnvironmentFile=/srv/app/.env
ExecStart=/srv/app/scripts/cloudflared-kv-updater.sh
Restart=always
RestartPreventExitStatus=78 79
RestartSec=20
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
